# tools/dharma-transcripts/launchd/com.jeffharris.dharma-ingestion.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jeffharris.dharma-ingestion</string>

  <key>WorkingDirectory</key>
  <string>/path/to/jeffharr.is</string>

  <key>ProgramArguments</key>
  <array>
    <string>/path/to/jeffharr.is/scripts/run-dharma-ingestion-keychain.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>BRENSILVER_INGEST_LIMIT</key>
    <string>20</string>
    <key>BRENSILVER_FEED_EVERY</key>
    <string>20</string>
    <key>BRENSILVER_AUTO_PUBLISH</key>
    <string>1</string>
    <key>BRENSILVER_MEDIA_BASE_URL</key>
    <string>https://jeffharr.is/dharma/brensilver/</string>
    <key>BRENSILVER_ARTWORK_BASE_URL</key>
    <string>https://jeffharr.is/dharma/brensilver/</string>
    <key>BRENSILVER_CHAPTERS_BASE_URL</key>
    <string>https://jeffharr.is/dharma/brensilver/</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>7200</integer>

  <key>StandardOutPath</key>
  <string>/path/to/jeffharr.is/.local-corpus/brensilver/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/path/to/jeffharr.is/.local-corpus/brensilver/logs/launchd.err.log</string>
</dict>
</plist>
