# tools/brensilver-transcripts/launchd/com.jeffharris.brensilver-transcripts.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jeffharris.brensilver-transcripts</string>

  <key>WorkingDirectory</key>
  <string>/Users/embergpt/code/flipper/brensilver/site</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/embergpt/code/flipper/brensilver/site/scripts/run-brensilver-ingestion.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>BRENSILVER_INGEST_LIMIT</key>
    <string>20</string>
    <key>BRENSILVER_FEED_EVERY</key>
    <string>20</string>
    <key>BRENSILVER_AUTO_PUBLISH</key>
    <string>1</string>
  </dict>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>StandardOutPath</key>
  <string>/Users/embergpt/code/flipper/brensilver/site/.local-corpus/brensilver/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/embergpt/code/flipper/brensilver/site/.local-corpus/brensilver/logs/launchd.err.log</string>
</dict>
</plist>
