# tools/brensilver-transcripts/launchd/com.jeffharris.brensilver-transcripts.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jeffharris.brensilver-transcripts</string>

  <key>WorkingDirectory</key>
  <string>/Users/embergpt/code/flipper/brensilver/site</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>PYTHONPATH=tools/brensilver-transcripts/src</string>
    <string>python3</string>
    <string>-m</string>
    <string>brensilver_transcripts.pipeline</string>
    <string>sync</string>
    <string>--limit</string>
    <string>3</string>
    <string>--update-qmd</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/embergpt/code/flipper/brensilver/site/.local-corpus/brensilver/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/embergpt/code/flipper/brensilver/site/.local-corpus/brensilver/logs/launchd.err.log</string>
</dict>
</plist>
